# services/plexpy/run
#!/bin/bash
exec /sbin/setuser abc /usr/bin/python /opt/plexpy/PlexPy.py --datadir /config 