# root/etc/services.d/tautulli/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /app/tautulli/PlexPy.py \
	--datadir /config
