# root/etc/services.d/plexpy/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /opt/plexpy/PlexPy.py \
	--datadir /config
